# templates/dokku-letsencrypt-auto-renew.service.j2
[Unit]
Description=Dokku letsencrypt:auto-renew

[Service]
Type=oneshot
ExecStart=/usr/bin/dokku letsencrypt:auto-renew
User=dokku
